# laravel/etc/supervisord/conf.d/laravel-worker.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
# command=php /プロジェクトまでのパス/artisan queue:work
command=php.exe artisan queue:work --sleep=3
autostart=true
autorestart=true
user=minatsuku
numprocs=2
redirect_stderr=true
# stdout_logfile=/プロジェクトまでのパス/worker.log
stdout_logfile=./etc/supervisord/tmp/worker.log
